# build/debian/postinst
#!/bin/sh -e

case "$1" in
	configure)
		if [ -x "/etc/init.d/sanewall" ]; then
			update-rc.d sanewall defaults >/dev/null
		fi
		;;

	abort-upgrade|abort-remove|abort-deconfigure)
		;;
esac
