# build/debian/preinst
#!/bin/sh -e

if [ -d /etc/firehol -a ! -d /etc/sanewall ]
then
	mv -f /etc/firehol /etc/sanewall
	echo
	echo
	echo "Changing your existing FireHOL configuration into a sanewall one"
	echo "configuration."
	echo "  /etc/firehol => /etc/sanewall"
	if [ -f /etc/sanewall/firehol.conf ]
	then
		mv -f /etc/sanewall/firehol.conf /etc/sanewall/sanewall.conf
		echo "  /etc/sanewall/firehol.conf => /etc/sanewall/sanewall.conf"
	fi
	if [ -f /etc/default/firehol ]
	then
		sed -e "s/FIREHOL/SANEWALL/g" /etc/default/firehol > /etc/default/sanewall
		chmod 644 /etc/default/sanewall
		rm -f /etc/default/firehol
		echo "  /etc/default/firehol => /etc/default/sanewall"
	fi
	echo
fi

if [ -f /etc/init.d/firehol ]
then
	rm -f /etc/init.d/firehol
	update-rc.d firehol remove
	echo
fi
