# debian/postrm
#!/bin/sh -e

case "$1" in
	remove)
		;;

	purge)
		update-rc.d -f firehol remove >/dev/null
		;;

	upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
		;;
esac
